# deploy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null           ; log to stdout so Cloud Run captures it
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
user=appuser

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ── FastAPI (uvicorn) ─────────────────────────────────────────────────────────
; Primary Cloud Run ingress port. Streamlit calls this via localhost.
[program:fastapi]
command=uvicorn src.api.main:app
    --host 0.0.0.0
    --port 8080
    --workers 2
    --log-level info
    --no-access-log
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
; Ensure FastAPI starts first so Streamlit's health check on startup succeeds.
priority=10

; ── Streamlit ─────────────────────────────────────────────────────────────────
; Internal port — not exposed by Cloud Run directly. The dashboard calls
; the FastAPI layer on localhost:8080.
[program:streamlit]
command=streamlit run src/app.py
    --server.port=8501
    --server.address=0.0.0.0
    --server.headless=true
    --server.enableCORS=false
    --server.enableXsrfProtection=false
directory=/app
autostart=true
autorestart=true
startretries=5
; Give FastAPI 10 seconds to be ready before Streamlit starts its health check.
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
